# etc/systemd/user/recordgps.service
[Unit]
Description=record GPS Data
After=gpsd.socket

[Service]
ExecStart=/usr/bin/python /home/pi/GPS/recordgps.py
